# ansible/roles/xbbs_coordinator/templates/xbbs-coordinator.service.j2
[Unit]
Description=xbbs coordinator

[Service]
User=xbbs
Group=xbbs
ExecStart=/bin/sh -xec '. '{{ xbbs_virtualenv|quote }}/bin/activate'; exec xbbs-coordinator'

[Install]
WantedBy=default.target
