# ansible/roles/xbbs_coordinator/templates/xbbs-coordinator.service.j2
[Unit]
Description=xbbs coordinator

[Service]
User=xbbs
Group=xbbs
ExecStart=/usr/local/bin/xbbs-coordinator
# TODO(arsen): FIX THIS
Environment=XBPS_ARCH=x86_64

[Install]
WantedBy=default.target
